# service/bufanalyzer.service
[Unit]
Description=Buffer Overflow Analyzer Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/bufanalyzer-daemon
Restart=always
RestartSec=5
User=user-12-31
Group=user-12-31
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=bufanalyzer

[Install]
WantedBy=multi-user.target
