# service/bufanalyzer.service.in
[Unit]
Description=Buffer Overflow Analyzer Daemon (@RUN_USER@)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/bufanalyzer-daemon
Restart=always
RestartSec=5
User=@RUN_USER@
Group=@RUN_USER@
WorkingDirectory=/var/lib/bufanalyzer
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bufanalyzer

[Install]
WantedBy=multi-user.target
